# files/etc/systemd/system/pocketforge-firstboot.service
[Unit]
Description=PocketForge OS First Boot Setup
After=multi-user.target
ConditionPathExists=!/var/lib/pocketforge-os/firstboot-complete

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/pocketforge-os/firstboot.sh

[Install]
WantedBy=multi-user.target
